# btc-predictor.service
[Unit]
Description=Bitcoin Price Predictor with ML and Firebase
Documentation=https://github.com/akunsamtiga/BTCPredictor
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=stcautotrade
Group=stcautotrade
WorkingDirectory=/home/stcautotrade/btc-predictor

# Environment variables - PATH HARUS INCLUDE VENV
Environment="PATH=/home/stcautotrade/btc-predictor/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/stcautotrade/btc-predictor"

# TensorFlow optimization untuk CPU-only
Environment="TF_CPP_MIN_LOG_LEVEL=2"
Environment="TF_ENABLE_ONEDNN_OPTS=0"
Environment="OMP_NUM_THREADS=2"
Environment="TF_NUM_INTEROP_THREADS=1"
Environment="TF_NUM_INTRAOP_THREADS=2"
Environment="KMP_AFFINITY=granularity=fine,compact,1,0"
Environment="KMP_BLOCKTIME=0"

# Load .env file (CRITICAL - must exist!)
EnvironmentFile=/home/stcautotrade/btc-predictor/.env

# Main command - Direct execution
ExecStart=/home/stcautotrade/btc-predictor/venv/bin/python3 /home/stcautotrade/btc-predictor/scheduler.py

# Pre-start check (optional but recommended)
ExecStartPre=/bin/bash -c 'test -f /home/stcautotrade/btc-predictor/.env || (echo "ERROR: .env file not found" && exit 1)'
ExecStartPre=/bin/bash -c 'test -f /home/stcautotrade/btc-predictor/venv/bin/python3 || (echo "ERROR: venv not found" && exit 1)'

# Post-start notification (optional)
# ExecStartPost=/bin/bash -c 'echo "Bitcoin Predictor started at $(date)" >> /home/stcautotrade/btc-predictor/logs/service.log'

# ALWAYS restart on failure
Restart=always
RestartSec=10

# Restart limits - Allow more attempts
StartLimitBurst=10
StartLimitIntervalSec=300

# Timeout settings
TimeoutStartSec=300
TimeoutStopSec=60

# Resource limits - Adjusted for VPS
MemoryMax=2G    
MemoryHigh=1.8G      
MemorySwapMax=2G     
CPUQuota=150%

# Nice level (lower priority than system services)
Nice=10

# Logging - systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc-predictor

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

# Writable paths - IMPORTANT!
ReadWritePaths=/home/stcautotrade/btc-predictor
ReadWritePaths=/home/stcautotrade/btc-predictor/logs
ReadWritePaths=/home/stcautotrade/btc-predictor/models
ReadWritePaths=/home/stcautotrade/btc-predictor/models_backup

# Kill behavior
KillMode=mixed
KillSignal=SIGTERM
SendSIGKILL=yes
FinalKillSignal=SIGKILL

[Install]
WantedBy=multi-user.target